# system/backend/fastapi-api.service
[Unit]
Description=Uvicorn (FastAPI) backend server
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/opt/app
ExecStart=/usr/bin/env uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
KillSignal=SIGTERM
TimeoutStopSec=5
PrivateTmp=true

[Install]
WantedBy=multi-user.target
